# bin/rc.d/alignak-reactionner
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: alignak_reactionner
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="alignak_reactionner"
rcvar="alignak_reactionner_enable"

alignak_reactionner_configfile="/usr/local/etc/alignak/daemons/reactionnerd.ini"
command="/usr/local/bin/alignak-reactionner"
command_interpreter="/usr/local/bin/python2.7"
command_args="-d -c ${alignak_reactionner_configfile} > /dev/null 2>&1"
pidfile="/var/run/alignak/reactionnerd.pid"

required_files="${alignak_reactionner_configfile}"

load_rc_config "${name}"

[ -z "${alignak_reactionner_enable}" ] && alignak_reactionner_enable="NO"

run_rc_command "$1"
