# bin/rc.d/alignak-receiver
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: alignak_receiver
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="alignak_receiver"
rcvar="alignak_receiver_enable"

alignak_receiver_configfile="/usr/local/etc/alignak/daemons/receiverd.ini"
command="/usr/local/bin/alignak-receiver"
command_interpreter="/usr/local/bin/python2.7"
command_args="-d -c ${alignak_receiver_configfile} > /dev/null 2>&1"
pidfile="/var/run/alignak/receiverd.pid"

required_files="${alignak_receiver_configfile}"

load_rc_config "${name}"

[ -z "${alignak_receiver_enable}" ] && alignak_receiver_enable="NO"

run_rc_command "$1"
