# bin/rc.d/alignak-poller
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: alignak_poller
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="alignak_poller"
rcvar="alignak_poller_enable"

alignak_poller_configfile="/usr/local/etc/alignak/daemons/pollerd.ini"
command="/usr/local/bin/alignak-poller"
command_interpreter="/usr/local/bin/python2.7"
command_args="-d -c ${alignak_poller_configfile} > /dev/null 2>&1"
pidfile="/var/run/alignak/pollerd.pid"

required_files="${alignak_poller_configfile}"

load_rc_config "${name}"

[ -z "${alignak_poller_enable}" ] && alignak_poller_enable="NO"

run_rc_command "$1"
