# bin/rc.d/alignak-scheduler
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: alignak_scheduler
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="alignak_scheduler"
rcvar="alignak_scheduler_enable"

alignak_scheduler_configfile="/usr/local/etc/alignak/daemons/schedulerd.ini"
command="/usr/local/bin/alignak-scheduler"
command_interpreter="/usr/local/bin/python2.7"
command_args="-d -c ${alignak_scheduler_configfile} > /dev/null 2>&1"
pidfile="/var/run/alignak/schedulerd.pid"

required_files="${alignak_scheduler_configfile}"

load_rc_config "${name}"

[ -z "${alignak_scheduler_enable}" ] && alignak_scheduler_enable="NO"

run_rc_command "$1"
