# bin/rc.d/alignak-broker
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: alignak_broker
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="alignak_broker"
rcvar="alignak_broker_enable"

alignak_broker_configfile="/usr/local/etc/alignak/daemons/brokerd.ini"
command="/usr/local/bin/alignak-broker"
command_interpreter="/usr/local/bin/python2.7"
command_args="-d -c ${alignak_broker_configfile} > /dev/null 2>&1"
pidfile="/var/run/alignak/brokerd.pid"

required_files="${alignak_broker_configfile}"

load_rc_config "${name}"

[ -z "${alignak_broker_enable}" ] && alignak_broker_enable="NO"

run_rc_command "$1"
