# bin/systemV/init.d/alignak-arbiter
#!/bin/sh
# Initially generated by sysd2v v0.3  --  http://www.trek.eu.org/devel/sysd2v
### BEGIN INIT INFO
# Provides:       alignak-arbiter
# Required-Start: $network $remote_fs $syslog
# Required-Stop:  $network $remote_fs $syslog
# Should-Start:
# Should-Stop:
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Description:    Alignak arbiter daemon instance
### END INIT INFO

# kFreeBSD do not accept scripts as interpreters, using #!/bin/sh and sourcing.
if [ true != "$INIT_D_SCRIPT_SOURCED" ] ; then
    set "$0" "$@"; INIT_D_SCRIPT_SOURCED=true . /lib/init/init-d-script
fi
set -a
ALIGNAK_CONFIGURATION_FILE=/usr/local/share/alignak/etc/alignak.ini
ALIGNAK_USER=alignak
ALIGNAK_GROUP=alignak
[ -r /etc/default/alignak ] && . /etc/default/alignak

# Alignak debs get installed inside "site-packages", but Debian's Python doesn't include it inside sys.path
# So let's find out which Python Alignak is going to use
PYTHONVER=$( /usr/bin/env python -V 2>&1 | cut -d' ' -f 2 | cut -d. -f 1,2 )
PYTHONPATH=/usr/local/lib/python$PYTHONVER/site-packages
set +a

NAME=arbiter-master
DESC=alignak-arbiter
DAEMON=/usr/local/bin/alignak-arbiter
PIDFILE=/var/run/alignak/arbiter-master.pid
LOGFILE=/var/log/alignak/arbiter-master.log
DAEMON_ARGS="-d -n ${NAME} -e "${ALIGNAK_CONFIGURATION_FILE}" -i ${PIDFILE} -l ${LOGFILE}"

do_start_cmd_override() {
  /bin/mkdir -p /var/log/alignak || return 2
  /bin/chown -R ${ALIGNAK_USER}:${ALIGNAK_GROUP} /var/log/alignak || return 2
  /bin/chmod 775 /var/log/alignak || return 2
  /bin/mkdir -p /var/run/alignak || return 2
  /bin/chown -R ${ALIGNAK_USER}:${ALIGNAK_GROUP} /var/run/alignak || return 2
  /bin/chmod 775 /var/run/alignak || return 2
  do_start_cmd || return $?
}

do_status_override() {
  status_of_proc -p "$PIDFILE" "$DAEMON" "$NAME" && return 0 || return $?
}

do_stop_cmd_override() {
  start-stop-daemon --stop --quiet --pidfile "$PIDFILE" --user "$ALIGNAK_USER" || return 2
}

do_reload() {
  log_daemon_msg "Reloading $DESC configuration files" "$NAME"
  start-stop-daemon --stop --signal HUP --quiet --pidfile "$PIDFILE" --user "$ALIGNAK_USER"
  log_end_msg $?
}
