# bin/systemV/init.d/alignak-webui
#!/bin/sh
# Generated by sysd2v v0.3  --  http://www.trek.eu.org/devel/sysd2v
# kFreeBSD do not accept scripts as interpreters, using #!/bin/sh and sourcing.
if [ true != "$INIT_D_SCRIPT_SOURCED" ] ; then
    set "$0" "$@"; INIT_D_SCRIPT_SOURCED=true . /lib/init/init-d-script
fi
### BEGIN INIT INFO
# Provides:       alignak-webui
# Required-Start: $remote_fs
# Required-Stop:  $remote_fs
# Should-Start:   $network
# Should-Stop:    $network
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Description:    uWSGI instance to serve Alignak WebUI
### END INIT INFO
set -a
ALIGNAK_WEBUI_CONFIGURATION_FILE=/usr/local/share/alignak-webui/etc/settings.cfg
ALIGNAK_WEBUI_BACKEND=http://127.0.0.1:5000
ALIGNAK_WEBUI_WS=http://127.0.0.1:7770
ALIGNAK_WEBUI_PID=/var/run//alignak-webui/alignak-webui.pid
ALIGNAK_WEBUI_LOG=/var/log/alignak-webui/alignak-webui.log
ALIGNAK_WEBUI_UWSGI_LOG=/var/log/alignak-webui/uwsgi.log
ALIGNAK_WEBUI_LOGGER_FILE=/usr/local/share/alignak-webui/etc/logging.json
ALIGNAK_WEBUI_HOST=127.0.0.1
ALIGNAK_WEBUI_PORT=5001
ALIGNAK_WEBUI_PROCESSES=4
ALIGNAK_USER=alignak
ALIGNAK_GROUP=alignak
[ -r /etc/default/alignak-webui ] && . /etc/default/alignak-webui

PYTHONVER=$( /usr/bin/env python -V 2>&1 | cut -d' ' -f 2 | cut -d. -f 1,2 )
PYTHONPATH=/usr/local/lib/python$PYTHONVER/site-packages
set +a

DESC=alignak-webui
NAME=alignak-webui
DAEMON=/usr/bin/uwsgi
DAEMON_ARGS="--master --enable-threads --plugin python --wsgi-file /usr/local/share/alignak-webui/etc/alignak-webui.wsgi --daemonize ${ALIGNAK_WEBUI_UWSGI_LOG} --pidfile ${ALIGNAK_WEBUI_PID} --logto /${ALIGNAK_WEBUI_LOG} --uid ${ALIGNAK_USER} --gid ${ALIGNAK_GROUP} --http-socket ${ALIGNAK_WEBUI_HOST}:${ALIGNAK_WEBUI_PORT} --processes ${ALIGNAK_WEBUI_PROCESSES} --ini /usr/local/share/alignak-webui/etc/uwsgi-service.ini"
PIDFILE=/var/run/alignak-webui/alignak-webui.pid

do_start_cmd_override ()
{
  /bin/mkdir -p /var/run/alignak-webui || return 2
  /bin/chown alignak:alignak /var/run/alignak-webui || return 2
  /bin/mkdir -p /var/log/alignak-webui || return 2
  /bin/chown alignak:alignak /var/log/alignak-webui || return 2
  do_start_cmd || return $?
}

do_stop_cmd_override() {
  start-stop-daemon --stop --signal INT --quiet --pidfile "$PIDFILE" --user "$ALIGNAK_USER" || return 2
}
